# aconfmgr/files/etc/udev/rules.d/pm-light.rules
# connect
SUBSYSTEM=="power_supply", ATTR{online}=="1", RUN+="/usr/bin/light -S 40"
