# aconfmgr/files/etc/systemd/system/lm_sensors.timer
[Timer]
OnBootSec=2sec

[Install]
WantedBy=timers.target
